# training/cloud/azure/files/usr/lib/systemd/system/nm-cloud-setup.service.d/10-rh-enable-for-azure.conf
[Service]
Environment="NM_CLOUD_SETUP_AZURE=yes"
